# api_to_mqtt.service
[Unit]
Description=API to MQTT Bridge Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User and Group inferred from your prompt 'pi@pbrscada'
User=pi
Group=pi

# WorkingDirectory inferred from your prompt '~/api_to_mqtt'
WorkingDirectory=/home/pi/api_to_mqtt

# Updated ExecStart with the specific path you provided
ExecStart=/nix/var/nix/profiles/default/bin/nix develop --extra-experimental-features 'nix-command flakes' --command python3 api_to_mqtt.py

# Restart policy to keep the service running if it crashes
Restart=always
RestartSec=10

# Environment variables if needed
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target